# crontab/log.sh
#!/bin/sh
cd /Users/oyo/Desktop/Node/crontab
echo $(date +%Y-%m-%d-%H-%M-%S)"\n" >> log.js 